# conf/supervisor/apps/confd.conf
[program:confd]
command=/usr/local/bin/confd -backend etcd -node http://%(COREOS_PRIVATE_IPV4):2379
