# conf/supervisor/apps/confd-local.conf
[program:confd]
command=/usr/local/bin/confd -backend env
