# conf/supervisor/apps/uwsgi.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/web.ini
